# conf/etc/crontab
*/5   *  *   *   *  root  bitcoin-cli gettxoutsetinfo

